# deploy/com.deepseek.exam-bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.deepseek.exam-bridge</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/env</string>
        <string>python3</string>
        <!-- EDIT: path to exam-bridge.py -->
        <string>/Users/tanjiawen/Infrastructure/references/DeepSeek-TUI/deploy/exam-bridge.py</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>DEEPSEEK_API_BASE</key>
        <string>http://127.0.0.1:7878</string>
        <key>BRIDGE_HOST</key>
        <string>127.0.0.1</string>
        <key>BRIDGE_PORT</key>
        <string>8888</string>
        <!-- EDIT: change this to a strong secret -->
        <key>BRIDGE_API_KEY</key>
        <string>bridge-secret-change-me</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>5</integer>

    <key>StandardOutPath</key>
    <string>/usr/local/var/log/exam-bridge.log</string>
    <key>StandardErrorPath</key>
    <string>/usr/local/var/log/exam-bridge.err</string>

    <key>WorkingDirectory</key>
    <string>/Users/tanjiawen</string>
</dict>
</plist>
